# systemd/friday-agent.service
[Unit]
Description=friday PoC Agent
After=network.target

[Service]
Type=simple
User=%i
ExecStart=/usr/bin/env python3 /opt/friday/poc/agent.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
